# reboot.timer
[Unit]
Description=Run reboot.service 2 hours after boot

[Timer]
OnBootSec=2h
Persistent=true

[Install]
WantedBy=timers.target
