# remmina-restart.service
[Unit]
Description=Restart Remmina service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart remmina
