# restart.timer
[Unit]
Description=Run restart service every five days

[Timer]
OnBootSec=5min
OnUnitActiveSec=5d
AccuracySec=1min

[Install]
WantedBy=timers.target
