# rhs_restart.timer
[Unit]
Description=Run restart service every day

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
AccuracySec=1min

[Install]
WantedBy=timers.target
